# praktikum-modul-1-e12-main/task-3/crontabs
*/2 * * * * /home/sakuya/tugas_3/automation.sh
*/10 * * * * /home/sakuya/tugas_3/download.sh
0 */2 * * * /home/sakuya/tugas_3/archive.sh
