# praktikum-modul-1-e12-main/task-4/crontabs
*/5 * * * * /home/sakuya/tugas_4/minute5_log.sh


0 * * * * /home/sakuya/tugas_4/agg_5min_to_hour.sh


0 * * * * /home/sakuya/tugas_4/uptime_monitor.sh

0 0 * * * /home/sakuya/tugas_4/cleanup_log.sh
